# sci-misc/boinc-server-project-eligius/files/boinc-server-eligius
#!/sbin/openrc-run

depend() {
	use dns net cron
}

check_baselayout() {
	if [ "${RC_VERSION:-0}" = "0" ]; then
		eerror "This script cannot be used for baselayout-1."
		return 1
	fi
}

start() {
	check_baselayout || return 1
	ebegin "Starting ${SVCNAME}"
	su -c "/usr/bin/python2 /var/lib/boinc-server/projects/eligius/bin/start" boincadm
	eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	su -c "/usr/bin/python2 /var/lib/boinc-server/projects/eligius/bin/stop" boincadm
	eend $?
}
