# app-arch/appimaged/files/appimaged-openrc
#!/sbin/openrc-run
# Tries to emulate the systemd user service but OpenRC way.

SERVER_PID_DIR="/var/run/appimaged/"
LOGFILE="/var/log/appimaged.log"
DEBUG=${DEBUG:=0}
APPIMAGED_ARGS=${APPIMAGED_ARGS:=""} # use -v for verbose

logmsg() {
	local msg="${1}"
	touch "${LOGFILE}"
	chown appimaged:appimaged "${LOGFILE}"
	if [[ "${ENABLE_LOGGING}" == "1" ]] ; then
		echo -e "${msg}" >> "${LOGFILE}"
	fi
}

get_appimage_users() {
	echo $(grep -P -e "^users" /etc/group | cut -f 4 -d ":" | tr "," " ")
}

start_servers() {
	mkdir -p "${SERVER_PID_DIR}"
	chown appimaged:appimaged "${SERVER_PID_DIR}"
	for user in $(get_appimage_users) ; do
		einfo "Starting appimaged for ${user}"
		local user_pid_file="${SERVER_PID_DIR}/${user}/${user}.pid"
		mkdir -p "$(dirname ${user_pid_file})"
		chown ${USER}:${USER} "$(dirname ${user_pid_file})"
		local extra_args
		if [[ "${DEBUG}" == "1" ]] ; then
			mkdir -p /var/log/appimaged/logs
			chown appimaged:appimaged /var/log/appimaged/logs
			chmod go+w /var/log/appimaged/logs
			ewarn "Logging stored in /var/log/appimaged/logs/${user}.log.  Remove when done for privacy and security reasons."
			extra_args="
				--stdout /var/log/appimaged/logs/${user}.log
				--stderr /var/log/appimaged/logs/${user}.log
			"
		fi
		start-stop-daemon \
			--background \
			--start \
			--group ${user} \
			--user ${user} \
			--pidfile "${user_pid_file}" \
			--make-pidfile \
			${extra_args} \
			--exec /usr/bin/appimaged -- ${APPIMAGED_ARGS}
	done
}

stop_servers() {
	for user in $(get_appimage_users) ; do
		einfo "Stopping appimaged for ${user}"
		local user_pid_file="${SERVER_PID_DIR}/${user}/${user}.pid"
		start-stop-daemon \
			--stop \
			--pidfile "${user_pid_file}"
	done
}

start()
{
	local m="Starting appimaged servers"
	ebegin "${m}"
	logmsg "${m}"
	start_servers
}

stop() {
	local m="Stopping appimaged servers"
	ebegin "${m}"
	logmsg "${m}"
	stop_servers
}
