# sys-apps/coolercontrold/files/coolercontrold.service
[Unit]
Description=CoolerControl Daemon
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=simple
Environment="CC_LOG=INFO"
ExecStart=/usr/bin/coolercontrold
Restart=always
RestartSec=1
TimeoutStopSec=10
User=coolercontrold
Group=coolercontrold

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
PrivateDevices=no          # needed for hwmon access
ProtectKernelTunables=yes
ProtectKernelModules=no    # may need modules
ProtectControlGroups=yes

ReadWritePaths=/etc/coolercontrol /var/lib/coolercontrol /run

[Install]
WantedBy=multi-user.target
