# sys-apps/coolercontrold/files/coolercontrold
#!/sbin/openrc-run
# This file is part of the coolercontrold package.

description="CoolerControl Daemon"

COOLERCONTROL_LOG="INFO"

pidfile="/run/${RC_SVCNAME}.pid"
logfile="/var/log/coolercontrold.log"

depend()
{
        need localmount
        after bootmisc coolercontrol-liqctld net
}

start_pre() {
	checkpath -d -m 0755 -o coolercontrold:coolercontrold /var/lib/coolercontrol
	checkpath -f -m 0644 -o coolercontrold:coolercontrold "${logfile}"

	# Ensure config directory is writable by the daemon
	checkpath -d -m 0755 -o coolercontrold:coolercontrold /etc/coolercontrol
}

start()
{
	ebegin "Starting coolercontrold"
		start-stop-daemon \
			--start \
			--user coolercontrold \
			--group coolercontrold \
			--pidfile "${pidfile}" \
			--background \
			--exec /usr/bin/coolercontrold
	eend $?
}

stop()
{
        local retval
	ebegin "Stopping coolercontrold"
		start-stop-daemon \
			--stop \
			--retry TERM/10/KILL/5 \
			--pidfile "${pidfile}"
	        retval=$?
        eend ${retval}
        return ${retval}
}
