# net-p2p/ckpool/files/ckdb
#!/sbin/openrc-run
# Copyright 1999-2017 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

checkconfig() {
	if [ ! -f ${CONFIG} ] ; then
		eerror "Config file \"${CONFIG}\" does not exist. Create it prior launching the daemon."
		return 1
	fi
}

depend() {
	after net
}

start() {
	checkconfig || return 1
	ebegin "Starting ckdb"
	start-stop-daemon --start --quiet --make-pidfile --pidfile ${CKDB_PIDFILE} --background --exec /usr/bin/ckdb -- ${CKDB_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping ckdb"
	start-stop-daemon --stop --quiet --pidfile ${CKDB_PIDFILE} --exec /usr/bin/ckdb
	eend $?
}
