# net-p2p/ckpool/files/ckpool
#!/sbin/openrc-run
# Copyright 1999-2017 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

env_check() {
	# Make sure the configuration is sane
	: ${CKPOOL_CONFIG:="/etc/ckpool/ckpool.conf"}
	return 0
}

checkconfig() {
	if [ ! -f ${CONFIG} ] ; then
		eerror "Config file \"${CONFIG}\" does not exist. Create it prior launching the daemon."
		return 1
	fi
}

depend() {
	after net
}

start() {
	checkconfig || return 1
	ebegin "Starting ckpool"
	start-stop-daemon --start --quiet --make-pidfile --pidfile ${CKPOOL_PIDFILE} --background --exec /usr/bin/ckpool -- -c "${CKPOOL_CONFIG}" ${CKPOOL_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping ckpool"
	start-stop-daemon --stop --quiet --pidfile ${CKPOOL_PIDFILE} --exec /usr/bin/ckpool
	eend $?
}
