# media-sound/tizonia/files/tizrmd
#!/sbin/openrc-run
# Copyright (c) 2007-2015 The OpenRC Authors.
# See the Authors file at the top-level directory of this distribution and
# https://github.com/OpenRC/openrc/blob/master/AUTHORS
#
# This file is part of OpenRC. It is subject to the license terms in
# the LICENSE file found in the top-level directory of this
# distribution and at https://github.com/OpenRC/openrc/blob/master/LICENSE
# This file may not be copied, modified, propagated, or distributed
# except according to the terms contained in the LICENSE file.

description="Tizonia's resource manager daemon"
PIDDIR=/run/tizonia
PIDFILE=$PIDDIR/$SVCNAME.pid
PIDFILE097=/var/run/$SVCNAME.pid

depend()
{
	need dbus
}

start()
{
	ebegin "Starting ${description}"

	local _retval=$?
	start-stop-daemon \
		--start \
		--group tizonia \
		--user tizonia \
		--quiet \
		--pidfile ${PIDFILE} \
		--exec /usr/bin/tizrmd -- --pid-file ${PIDFILE}
	eend $?
}

stop() {
	ebegin "Stopping ${description}"
	if [ ! -e "${PIDFILE}" ] && [ -e ${PIDFILE097} ]; then
		# Try stopping the pid file used by <0.9.7
		start-stop-daemon --stop --quiet --pidfile "${PIDFILE097}"
		ret=$?
		[ ${ret} = 0 ] && [ -e "${PIDFILE097}" ] && rm "${PIDFILE097}"
		eend ${ret}
	else
		start-stop-daemon --stop --quiet --pidfile "${PIDFILE}"
		eend $?
	fi
}
