# sys-kernel/ot-sources/files/ot-kernel-iosched.openrc
#!/sbin/openrc-run
# This file is part of the ot-sources package.

description="Sets I/O scheduler of each drive individually."

depend()
{
        need localmount
}

start()
{
	/usr/bin/ot-kernel-iosched.sh
}
