# app-misc/ollama/files/ollama.openrc
#!/sbin/openrc-run

description="Ollama Service"
real_command="/usr/bin/ollama"
real_command_args="serve"
command_user="ollama"
command_group="ollama"
command_background="yes"
pidfile="/run/ollama.pid"
log="/var/log/ollama/ollama.log"

# AC = Access Control
# CE = Code Execution
# ID = Information Disclosure
# DoS = Denial of Service
# DT = Data Tampering

# CWE-653 mitigation with Firejail sandboxing
# Use the Firejail to mitigate against trojanized LLMs with .so files or malicous payloads.

_OLLAMA_DEBUG=${_OLLAMA_DEBUG:-"0"}

OLLAMA_BACKEND=${OLLAMA_BACKEND:-"@OLLAMA_BACKEND@"}
OLLAMA_CONTEXT_LENGTH=${OLLAMA_CONTEXT_LENGTH:-@OLLAMA_CONTEXT_LENGTH@}
OLLAMA_FIREJAIL=${OLLAMA_FIREJAIL:-@OLLAMA_FIREJAIL@}
OLLAMA_FLASH_ATTENTION=${OLLAMA_FLASH_ATTENTION:-@OLLAMA_FLASH_ATTENTION@}
OLLAMA_KV_CACHE_TYPE=${OLLAMA_KV_CACHE_TYPE:-"@OLLAMA_KV_CACHE_TYPE@"}
OLLAMA_VULKAN=${OLLAMA_VULKAN:-"@OLLAMA_VULKAN@"}

# CWE-401, DoS
# Kill process on OOM for FAFO users
# -1000 never killed
# 0 default
# <500 mid-back of the OOM kill list
# >=500 mid-front of the OOM kill list
# 1000 front of the line OOM kill
OLLAMA_OOM=${OLLAMA_OOM:-500}

# Ollama allows cross-origin requests from 127.0.0.1 and 0.0.0.0 by default.
# Additional origins can be configured with OLLAMA_ORIGINS.
# export OLLAMA_ORIGINS="<ip>"

bind_mount_file() {
	local path="${1}"
	local dir
}

start() {
	ebegin "Starting $description"
	exec >> >(logger -t "$RC_SVCNAME Start daemon" -p "daemon.info")

	local uid=$(grep "^ollama:" "/etc/passwd" | cut -f 3 -d ":")
	local gid=$(grep "^ollama:" "/etc/passwd" | cut -f 4 -d ":")

	local args=(
		-e LD_LIBRARY_PATH="@LD_LIBRARY_PATH@"
	)

	# See also https://github.com/ollama/ollama/blob/v0.30.6/envconfig/config.go
	# Editor blocked to limit implied blast radius of saving or accessing home folder.
	local L=(
		"CUDA_VISIBLE_DEVICES"
		#"EDITOR"
		"GGML_VK_VISIBLE_DEVICES"
		"GPU_DEVICE_ORDINAL"
		"HIP_VISIBLE_DEVICES"
		"HSA_OVERRIDE_GFX_VERSION"
		"HTTP_PROXY"
		"http_proxy"
		"HTTPS_PROXY"
		"https_proxy"
		"LLAMA_ARG_FIT"
		"LLAMA_ARG_FIT_TARGET"
		"NO_PROXY"
		"no_proxy"
		"OLLAMA_DEBUG"
		"OLLAMA_DEBUG_LOG_REQUESTS"
		"OLLAMA_CONTEXT_LENGTH"
		#"OLLAMA_EDITOR"
		"OLLAMA_FLASH_ATTENTION"
		"OLLAMA_GPU_OVERHEAD"
		"OLLAMA_GO_TEMPLATE"
		"OLLAMA_IGPU_ENABLE"
		"OLLAMA_KV_CACHE_TYPE"
		"OLLAMA_HOST"
		"OLLAMA_KEEP_ALIVE"
		"OLLAMA_LLM_LIBRARY"
		"OLLAMA_LOAD_TIMEOUT"
		"OLLAMA_MAX_LOADED_MODELS"
		"OLLAMA_MAX_TRANSFER_STREAMS"
		"OLLAMA_MAX_QUEUE"
		"OLLAMA_MODELS"
		"OLLAMA_NO_CLOUD"
		"OLLAMA_NOHISTORY"
		"OLLAMA_NOPRUNE"
		"OLLAMA_NUM_PARALLEL"
		"OLLAMA_REMOTES"
		"OLLAMA_ORIGINS"
		"OLLAMA_SCHED_SPREAD"
		"OLLAMA_VULKAN"
		"ROCR_VISIBLE_DEVICES"
		#"VISUAL"
	)

	local x
	for x in "${L[@]}" ; do
		if [[ -n "${!x}" ]] ; then
			args+=(
				-e ${x}="${!x}"
			)
		fi
	done

	if [[ "${_OLLAMA_DEBUG}" == "1" ]] ; then
		:
	else
		args+=(
			--background
			--stdout "$log"
			--stderr "$log"
		)
	fi

	args+=(
		--exec $real_command -- $real_command_args
	)

	start-stop-daemon \
		--start \
		--make-pidfile \
		--pidfile "$pidfile" \
		--user "$command_user" \
		--group "$command_group" \
		${args[@]}

	eend $?
}

stop() {
	ebegin "Stopping $description"
	exec >> >(logger -t "$RC_SVCNAME Stop daemon" -p daemon.info)

#	start-stop-daemon --stop --pidfile "$pidfile"

	# The start-stop-daemon line is bugged
	pkill -9 -G ollama

	eend $?
}
