# app-misc/ollama/files/ollama
#!/sbin/openrc-run

description="Ollama Service"
real_command="/usr/bin/ollama"
real_command_args="serve"
command_user="ollama"
command_group="ollama"
command_background="yes"
pidfile="/run/ollama.pid"
log="/var/log/ollama/ollama.log"

# AC = Access Control
# CE = Code Execution
# ID = Information Disclosure
# DoS = Denial of Service
# DT = Data Tampering

# CWE-653 mitigation
OLLAMA_CHROOT=${OLLAMA_CHROOT:-@OLLAMA_CHROOT@}

OLLAMA_CONTEXT_LENGTH=${OLLAMA_CONTEXT_LENGTH:-@OLLAMA_CONTEXT_LENGTH@}

OLLAMA_FLASH_ATTENTION=${OLLAMA_FLASH_ATTENTION:-@OLLAMA_FLASH_ATTENTION@}

OLLAMA_KV_CACHE_TYPE=${OLLAMA_KV_CACHE_TYPE:-"@OLLAMA_KV_CACHE_TYPE@"}

# CWE-653 mitigation
OLLAMA_SANDBOX_PROVIDER=${OLLAMA_SANDBOX_PROVIDER:-"@OLLAMA_SANDBOX_PROVIDER@"}

OLLAMA_BACKEND=${OLLAMA_BACKEND:-"@OLLAMA_BACKEND@"}
_OLLAMA_DEBUG=${_OLLAMA_DEBUG:-"0"}
OLLAMA_VULKAN=${OLLAMA_VULKAN:-"@OLLAMA_VULKAN@"}

# 0 back of the line OOM kill list
# <500 back of the OOM kill list
# >=500 front of the OOM kill list
# 1000 front of the line OOM kill
OLLAMA_OOM=${OLLAMA_OOM:-500}

# Ollama allows cross-origin requests from 127.0.0.1 and 0.0.0.0 by default.
# Additional origins can be configured with OLLAMA_ORIGINS.
# export OLLAMA_ORIGINS="<ip>"

bind_mount_file() {
	local path="${1}"
	local dir
}

prepare_chroot() {
	[[ -e "/var/lib/ollama" ]] || exit 1
	[[ -e "/var/lib/ollama" ]] || return
	pushd "/var/lib/ollama" >/dev/null 2>&1 || exit 1
		paths=(
			"lib64"
			"bin"
			"usr"
			"var"
			"etc"
			"tmp"
		)
		rm -rf ${paths[@]}

		umount -l "proc" >/dev/null 2>&1
		umount -l "sys" >/dev/null 2>&1
		umount -l "dev" >/dev/null 2>&1
		umount -l "run" >/dev/null 2>&1

		mkdir -p "proc"
		mkdir -p "sys"
		mkdir -p "dev"
		mkdir -p "run"

		local gcc_slot=$(ldd "/usr/lib/ollama/bin/ollama" \
			| grep "libstdc++.so" \
			| grep -E -o -e "/[0-9]+/" \
			| sed -e "s|/||g")
		local chost=$(cat "/var/db/pkg/app-misc/ollama"*"/CHOST")

		mkdir -p "var/log/ollama"
		mkdir -p "var/lib/ollama"
		touch "var/log/ollama/ollama.log"
		chown -R "ollama:ollama" "var/log/ollama"
		chown "ollama:ollama" "var/lib/ollama"
		chmod 0755 "var/lib"
		chmod 0700 "var/lib/ollama"
		chmod 0755 "/var/lib/ollama"
		mkdir -p "tmp"
		chmod 1777 "tmp"

		if [[ "${OLLAMA_SANDBOX_PROVIDER}" != "bwrap" ]] ; then
			mount --types "proc" "/proc" "proc"
			mount --rbind "/sys" "sys"
			mount --make-rslave "sys"
			mount --rbind "/dev" "dev"
			mount --make-rslave "dev"
			mount --bind "/run" "run"
			mount --make-slave "run"
		fi

		local paths=()
		paths=(
			"/etc/ld.so.cache"
			"/etc/ld.so.preload"
			"/bin/bash;libs"
			"/bin/sh"
			"/usr/bin/ollama"
			"/usr/lib/ollama"
			"/usr/lib/ollama;libs-recursive"
			"/bin/mkdir"
		)

		if [[ "${OLLAMA_SANDBOX_PROVIDER}" == "bwrap" ]] ; then
			paths+=(
				"/usr/bin/bwrap;libs"
			)
		elif [[ "${OLLAMA_SANDBOX_PROVIDER}" == "sandbox" ]] ; then
			paths+=(
				"/etc/sandbox.conf"
				"/usr/share/sandbox"
				"/usr/lib64/libsandbox.so"
				"/usr/bin/sandbox"
				"/usr/bin/unshare"
			)
		fi

	# For debugging only
		if [[ "${_OLLAMA_DEBUG}" == "1" ]] ; then
			paths+=(
				"/bin/ls"
				"/bin/cat"
				"/bin/grep"
				"/bin/echo"
				"/bin/pwd"
				"/usr/bin/top;libs"
				"/bin/ps;libs"
				"/usr/bin/id"
				"/usr/bin/ping"
			)
		fi

		paths+=(
			"/bin/getcwd;libs"
		)

		if [[ "${OLLAMA_BACKEND}" == "rocm" ]] ; then
			mkdir -p "opt/rocm"
			mount -o ro --bind "/opt/rocm" "opt/rocm"

			paths+=(
				"/opt/rocm;libs-recursive"
			)
		elif [[ "${OLLAMA_BACKEND}" == "cuda" ]] ; then
			mkdir -p "opt/cuda"
			mount -o ro --bind "/opt/cuda" "opt/cuda"
			paths+=(
				"/opt/cuda;libs-recursive"
			)
		fi

		local path
		for path in ${paths[@]} ; do
			if [[ "${path}" =~ ";libs-recursive" ]] ; then
	# Copy ldd libs recursively without copying folder
				local src="${path%;*}"
				IFS=$'\n'
				local pair
				for pair in $(ldd $(find "${src}" -executable -type f)) ; do
					[[ "${pair}" =~ "=>" ]] || continue
					local lib_path=$(echo "${pair}" \
						| cut -f 2 -d ">" \
						| cut -f 1 -d '(' \
						| sed -r -e "s|^[ ]+||g" -e "s|[ ]+$||g")
					local dir_path=$(dirname "${lib_path}")
					dir_path="${dir_path:1}"
					[[ -z "${dir_path}" ]] && continue
					[[ "${dir_path}" =~ "opt/rocm" ]] && continue
					[[ "${dir_path}" =~ "opt/cuda" ]] && continue
					mkdir -p "${dir_path}"
					cp -a "${lib_path}" "${dir_path}"
					if [[ -L "${lib_path}" ]] ; then
						local path2=$(realpath "${lib_path}")
						if [[ -f "${path2}" ]] ; then
							cp -a "${path2}" "${dir_path}"
						fi
					fi
				done

				local row
				for row in $(ldd $(find "${src}" -executable -type f)) ; do
					[[ "${pair}" =~ "(" ]] || continue
					[[ "${pair}" =~ "=>" ]] && continue
					local lib_path=$(echo "${row}" \
						| cut -f 1 -d "(" \
						| sed -r -e "s|^[[:space:]]+||g" -e "s|[[:space:]]+$||g")
					local dir_path=$(dirname "${lib_path}")
					dir_path="${dir_path:1}"
					if [[ -e "${lib_path}" ]] ; then
						mkdir -p "${dir_path}"
						cp -a "${lib_path}" "${dir_path}"
					fi
				done

				IFS=$' \t\n'
			elif [[ "${path}" =~ ";libs" ]] ; then
	# Direct copy with copy ldd libs
				local src="${path%;*}"

	# Direct copy
				local dest
				dest=$(dirname "${src}")
				dest="${dest:1}"
				mkdir -p "${dest}"
				cp -a "${src}" "${dest}"

	# Copy libs depends
				IFS=$'\n'
				local pair
				for pair in $(ldd "${src}") ; do
					[[ "${pair}" =~ "=>" ]] || continue
					local lib_path=$(echo "${pair}" \
						| cut -f 2 -d ">" \
						| cut -f 1 -d '(' \
						| sed -r -e "s|^[ ]+||g" -e "s|[ ]+$||g")
					local dir_path=$(dirname "${lib_path}")
					dir_path="${dir_path:1}"
					[[ -z "${dir_path}" ]] && continue
					[[ "${dir_path}" =~ "opt/rocm" ]] && continue
					[[ "${dir_path}" =~ "opt/cuda" ]] && continue
					mkdir -p "${dir_path}"
					cp -a "${lib_path}" "${dir_path}"
					if [[ -L "${lib_path}" ]] ; then
						local path2=$(realpath "${lib_path}")
						if [[ -f "${path2}" ]] ; then
							cp -a "${path2}" "${dir_path}"
						fi
					fi
				done

				local row
				for row in $(ldd $(find "${src}" -executable -type f)) ; do
					[[ "${pair}" =~ "(" ]] || continue
					[[ "${pair}" =~ "=>" ]] && continue
					local lib_path=$(echo "${row}" \
						| cut -f 1 -d "(" \
						| sed -r -e "s|^[[:space:]]+||g" -e "s|[[:space:]]+$||g")
					local dir_path=$(dirname "${lib_path}")
					dir_path="${dir_path:1}"
					if [[ -e "${lib_path}" ]] ; then
						mkdir -p "${dir_path}"
						cp -a "${lib_path}" "${dir_path}"
					fi
				done

				IFS=$' \t\n'
			else
	# Direct copy
				local dest
				dest=$(dirname "${path}")
				dest="${dest:1}"
				mkdir -p "${dest}"
				cp -a "${path}" "${dest}"
			fi
		done

	popd >/dev/null 2>&1
}

stop_chroot() {
	[[ -e "/var/lib/ollama" ]] || exit 1
	[[ -e "/var/lib/ollama" ]] || return
	pushd "/var/lib/ollama" >/dev/null 2>&1 || exit 1
		if [[ "${OLLAMA_SANDBOX_PROVIDER}" == "sandbox" ]] ; then
			umount -l "proc" >/dev/null 2>&1
			umount -l "sys" >/dev/null 2>&1
			umount -l "dev" >/dev/null 2>&1
			umount -l "run" >/dev/null 2>&1
		fi
		if [[ "${OLLAMA_BACKEND}" == "rocm" ]] ; then
			umount -l "opt/rocm"
		elif [[ "${OLLAMA_BACKEND}" == "cuda" ]] ; then
			umount -l "opt/cuda"
		fi
	popd >/dev/null 2>&1
}

start() {
	if [[ "${OLLAMA_CHROOT}" == "1" ]] ; then
		prepare_chroot
	fi

	ebegin "Starting $description"
	exec >> >(logger -t "$RC_SVCNAME Start daemon" -p "daemon.info")

	local args=()

	local uid=$(grep "^ollama:" "/etc/passwd" | cut -f 3 -d ":")
	local gid=$(grep "^ollama:" "/etc/passwd" | cut -f 4 -d ":")

	args+=(
#		-e OLLAMA_UID=${uid}
#		-e OLLAMA_GID=${gid}
		-e LD_LIBRARY_PATH="@LD_LIBRARY_PATH@"
	)

	if [[ "${OLLAMA_SANDBOX_PROVIDER}" == "sandbox" ]] ; then
		: # Broken
	elif [[ "${OLLAMA_SANDBOX_PROVIDER}" == "bwrap" ]] ; then
		: # Duplicate functionality
	elif [[ "${OLLAMA_CHROOT}" == "1" ]] ; then
		args+=(
			--chroot "/var/lib/ollama"
		)
	fi

	if [[ "${_OLLAMA_DEBUG}" == "1" ]] ; then
		:
	else
		args+=(
			--background
			--stdout "$log"
			--stderr "$log"
		)
	fi

	if [[ "${OLLAMA_SANDBOX_PROVIDER}" == "sandbox" ]] ; then
	# CWE-653 mitigation
		sandbox_command="/usr/bin/sandbox"
		sandbox_args=""
		args+=(
	# CWE-653 mitigation, ID
			-e SANDBOX_DENY="/proc:/dev/input:/root:/home:/dev/mem:/dev/kmem:/proc/kcore:/mnt:/etc/passwd:/etc/groups:${OLLAMA_SANDBOX_BLACKLISTED_PATHS}"
			--exec $sandbox_command -- $sandbox_args $real_command $real_command_args
		)
	else
		args+=(
			--exec $real_command -- $real_command_args
		)
	fi


	start-stop-daemon \
		--start \
		--make-pidfile \
		--pidfile "$pidfile" \
		--user "$command_user" \
		--group "$command_group" \
		${args[@]}

	# CWE-401, DoS
	# Kill on OOM
	sleep 1
	if choom --version >/dev/null 2>&1 ; then
		local pid
		for pid in $(pgrep -U ollama) ; do
			choom -p ${pid} -n ${OLLAMA_OOM}
		done
	fi

	eend $?
}

stop() {
	ebegin "Stopping $description"
	if [[ "${OLLAMA_CHROOT}" == "1" ]] ; then
		stop_chroot
	fi
	exec >> >(logger -t "$RC_SVCNAME Stop daemon" -p daemon.info)

#	start-stop-daemon --stop --pidfile "$pidfile"

	# The start-stop-daemon line is bugged
	pkill -9 -G ollama

	eend $?
}
