# media-radio/codecserver/files/init.d/codecserver
#!/sbin/openrc-run

description="Runs codecserver."

depend() {
	need localmount
}

start()
{
	einfo "Starting codecserver"
	if start-stop-daemon \
		--background \
		--make-pidfile \
		--pidfile /run/codecserver.pid \
		--start \
		--user codecserver:codecserver \
		/usr/bin/codecserver
	then
		eend 0
	else
		eend 1
	fi
}

stop() {
	start-stop-daemon \
		--stop \
		--pidfile /run/codecserver.pid
	eend 0
}
